# target/linux/ath25/base-files/etc/uci-defaults/01_leds
#!/bin/sh
# Copyright 2012 OpenWrt.org
#

. /lib/functions/uci-defaults.sh

ucidef_set_led_netdev "wlan" "wlan" "wlan" "wlan0"
# "ignored" so that the init.d script doesn't change it.
ucidef_set_led_default "etactica" "etactica" "globe" "ignored"

ucidef_commit_leds

exit 0
